# queueme-celery.service
[Unit]
Description=QueueMe Celery Worker
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/path/to/queueme
Environment="PATH=/path/to/queueme/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="DJANGO_SETTINGS_MODULE=queueme.settings.production"
ExecStart=/path/to/queueme/venv/bin/celery -A queueme worker --loglevel=info
Restart=on-failure
RestartSec=5
SyslogIdentifier=queueme-celery

[Install]
WantedBy=multi-user.target 